# cloudflared.service
[Unit]
Description=Cloudflare Tunnel – fashion-scraper public endpoint
After=network.target fashion-scraper.service
Requires=fashion-scraper.service

[Service]
Type=simple
User=parimal
ExecStart=/usr/local/bin/cloudflared tunnel --url http://localhost:8000 --no-autoupdate
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

